# cron.d/speedtest.sh
*/10 * * * * zabbix /etc/zabbix/scripts/speedtest.sh -a
